# chef/site-cookbooks/snow/templates/default/workers/upstart/rippleout.conf.erb
start on startup
env NAME=snow-rippleout
env NODE_ENV=<%= node.chef_environment %>
env DEBUG=*
env HOME=/home/ubuntu
chdir <%= node[:snow][:workers][:app_directory] %>/current/workers

script
    echo $$ > /var/run/$NAME.pid
    sudo -Eu ubuntu node bin/rippleout >> /var/log/$NAME.log 2>&1
end script

post-stop script
   rm -f /var/run/$NAME.pid
end script
