# chef/site-cookbooks/snow/templates/default/insight/upstart/insightservice.conf.erb
start on startup
console log
env NAME=<%= @service %>
env HOME=/home/ubuntu

script
    echo $$ > $NAME.pid
    sudo -Eu ubuntu /usr/bin/$NAME.sh
end script

post-stop script
   rm -f /var/run/$NAME.pid
end script
